# Resources/LaunchDaemons/com.zerofs.manager.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zerofs.manager.helper</string>
  <key>BundleProgram</key>
  <string>Contents/MacOS/ZeroFSPrivilegedHelper</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>MachServices</key>
  <dict>
    <key>com.zerofs.manager.helper</key>
    <true/>
  </dict>
  <key>AssociatedBundleIdentifiers</key>
  <array>
    <string>com.zerofs.manager</string>
  </array>
</dict>
</plist>
